# scripts/com.jpc.check-agent-drive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.jpc.check-agent-drive.plist
  ═══════════════════════════════════════════════════════════
  macOS Launch Agent — runs check_mount.sh at login/boot.
  Logs to the external drive if mounted, internal /tmp if not.

  INSTALL (run once on the Mac Mini):
    cp scripts/com.jpc.check-agent-drive.plist \
       ~/Library/LaunchAgents/com.jpc.check-agent-drive.plist

    launchctl load ~/Library/LaunchAgents/com.jpc.check-agent-drive.plist

  VERIFY:
    launchctl list | grep jpc

  UNINSTALL:
    launchctl unload ~/Library/LaunchAgents/com.jpc.check-agent-drive.plist
    rm ~/Library/LaunchAgents/com.jpc.check-agent-drive.plist
  ═══════════════════════════════════════════════════════════
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jpc.check-agent-drive</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/MAC_MINI_1TB/LegionForge/scripts/check_mount.sh</string>
        <string>--wait</string>
        <string>--create-dirs</string>
    </array>

    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it crashes (useful during boot race conditions) -->
    <key>KeepAlive</key>
    <false/>

    <!-- Wait 10s after login before running — gives drive time to mount -->
    <key>StartInterval</key>
    <integer>10</integer>

    <!-- Logs — written to external drive once mounted -->
    <key>StandardOutPath</key>
    <string>/tmp/jpc-mount-check.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/jpc-mount-check.err</string>

    <!-- Only run as the current user, not root -->
    <key>UserName</key>
    <string>jp</string>
</dict>
</plist>
